# templates/openwrt-getdomains.j2
#!/bin/sh /etc/rc.common

START=99
CACHE_DIR="/etc/routewolf/lists"
TMP_DNSMASQ_DIR="/tmp/dnsmasq.d"
TMP_LIST_DIR="/tmp/lst"

DOMAINS="{{ domain_list_url | default('https://raw.githubusercontent.com/dagmagnat/RouteWolf/main/lists/domains-dnsmasq-nfset.lst') }}"
IPV4="{{ ipv4_list_url | default('https://raw.githubusercontent.com/dagmagnat/RouteWolf/main/lists/ipv4.lst') }}"

start () {
    mkdir -p "$CACHE_DIR" "$TMP_DNSMASQ_DIR" "$TMP_LIST_DIR"

    [ -s "$CACHE_DIR/domains.lst" ] && cp "$CACHE_DIR/domains.lst" "$TMP_DNSMASQ_DIR/domains.lst"
    [ -s "$CACHE_DIR/ipv4.lst" ] && cp "$CACHE_DIR/ipv4.lst" "$TMP_LIST_DIR/ipv4.lst"

    if [ -n "$DOMAINS" ] && curl -L -f --connect-timeout 10 --retry 3 "$DOMAINS" --output "$TMP_DNSMASQ_DIR/domains.lst.new"; then
        if dnsmasq --conf-file="$TMP_DNSMASQ_DIR/domains.lst.new" --test 2>&1 | grep -q "syntax check OK"; then
            mv "$TMP_DNSMASQ_DIR/domains.lst.new" "$TMP_DNSMASQ_DIR/domains.lst"
            cp "$TMP_DNSMASQ_DIR/domains.lst" "$CACHE_DIR/domains.lst"
            /etc/init.d/dnsmasq restart
        else
            rm -f "$TMP_DNSMASQ_DIR/domains.lst.new"
            echo "Downloaded domain list is invalid; cached list was kept"
        fi
    fi

    {% if list_ip %}
    if [ -n "$IPV4" ] && curl -L -f --connect-timeout 10 --retry 3 "$IPV4" --output "$TMP_LIST_DIR/ipv4.lst.new"; then
        mv "$TMP_LIST_DIR/ipv4.lst.new" "$TMP_LIST_DIR/ipv4.lst"
        cp "$TMP_LIST_DIR/ipv4.lst" "$CACHE_DIR/ipv4.lst"
        /etc/init.d/firewall restart
    fi
    {% endif %}
}
